# per-forum.conf
[program:per-forum]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/root/per-forum
autostart=true
autorestart=true
